# deploy/rctf/supervisord.conf
[unix_http_server]
file=/tmp/supervisor/supervisor.sock

[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisor/supervisord.pid

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor/supervisor.sock

[program:api]
command=/usr/local/bin/bun run /app/apps/api/dist/index.js
user=rctf
autostart=true
autorestart=true
# outlive the API's graceful shutdown
stopwaitsecs=35
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
# Disabling rotation so that supervisor isn't seeking the stdout/stderr handles
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;' -e stderr
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
